# smart-cane.service
[Unit]
Description=Smart Cane System for Blind Users
Documentation=https://github.com/AmmarRiaz123/Navicane
After=multi-user.target network.target

[Service]
Type=simple
# User will be set dynamically by install.sh
User=pi
Group=pi
# WorkingDirectory will be set dynamically by install.sh
WorkingDirectory=/home/pi/smart_cane
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/pi/smart_cane/main.py
Restart=always
RestartSec=10
# Logs will go to working directory
StandardOutput=append:/home/pi/smart_cane/smart_cane.log
StandardError=append:/home/pi/smart_cane/smart_cane.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
MemoryMax=2G
CPUQuota=80%

[Install]
WantedBy=multi-user.target